# docs/systemd/archon-pg-backup.timer
[Unit]
Description=Run Archon Postgres backup nightly at 03:30

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=archon-pg-backup.service

[Install]
WantedBy=timers.target
