# docs/systemd/archon-pg-backup.service
[Unit]
Description=Archon Postgres nightly backup

[Service]
Type=oneshot
EnvironmentFile=%h/.archon/.env
ExecStart=%h/archon/scripts/pg-backup.sh
StandardOutput=journal
StandardError=journal
